# config/tribefinder-auto-backup.service
[Unit]
Description=TribeFinder - Automatic Backups
After=network.target

[Service]
Type=oneshot
User=tribefinder
Group=tribefinder
WorkingDirectory=/home/tribefinder/TribeFinder
Environment=NODE_ENV=production
EnvironmentFile=/home/tribefinder/TribeFinder/.env
ExecStart=/usr/bin/node /home/tribefinder/TribeFinder/scripts/auto-backup.js
