# config/tribefinder-marketplace-expiry.timer
[Unit]
Description=TribeFinder - Marketplace Expiry (Daily)

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=900

[Install]
WantedBy=timers.target
