# config/tribefinder-auto-backup.timer
[Unit]
Description=TribeFinder - Automatic Backups (Hourly Check)

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
